# config/cs.service
[Unit]
Wants=network-online.target
After=network-online.target nss-lookup.target

[Service]
User=qrystal-cs
Restart=on-failure
Type=notify
ExecStart=/opt/qrystal/qrystal-cs -config /etc/qrystal/cs-config.yml
RuntimeDirectory=qrystal-cs
RuntimeDirectoryMode=0700
StateDirectory=qrystal-cs
StateDirectoryMode=0700
LogsDirectory=qrystal-cs
LogsDirectoryMode=0700
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
