# config/runner.service
[Unit]
Description=Run Qrystal Node using Runner.
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
ExecStart=qrystal-runner -config /etc/qrystal/runner-config.yml
WorkingDirectory=/opt/qrystal
Environment=QRYSTAL_PROFILE_PATH=/tmp/qrystal-profile.pprof

[Install]
WantedBy=multi-user.target
